# package/base-files/files/etc/crontabs/root
#*/20 * * * * /etc/ppp/check-lte.sh
SHELL=/bin/bash
5 */1 * * * killall ntpd; sleep $(expr $RANDOM / 100); ntpd -n -N -p ntp1.ha-vel.cz ntp2.ha-vel.cz -q && /usr/bin/fake-hwclock save
