# package/base-files/files/etc/crontabs/root
SHELL=/bin/bash
1 */1 * * * sleep $(expr $RANDOM / 100); ntpd -n -N -p ntp1.ha-vel.cz ntp2.ha-vel.cz -q 2 &1> /dev/null 
